# systemd-wrappers/ayatana-webmail.service
[Unit]
Description=Ayatana Notification service

[Service]
ExecStart=/usr/bin/ayatana-webmail
Restart=always

[Install]
WantedBy=multi.user.target

